# app/Third Party API/indopoints/.ebextensions/crontab
*/5 * * * * python /opt/python/bundle/2/app/ipScheduler.py
