# app/Third Party API/gopoints/.ebextensions/crontab
*/5 * * * * python /opt/python/bundle/2/app/gpScheduler.py
